# deploy/ubuntu/graffiti-student.service
[Unit]
Description=Graffiti Student Local API
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=graffiti
Group=graffiti
WorkingDirectory=/srv/graffiti-student/app
EnvironmentFile=/etc/graffiti-student.env
Environment=PYTHONUNBUFFERED=1
ExecStart=/srv/graffiti-student/venv/bin/uvicorn deploy.local_api:app --host 127.0.0.1 --port 8000 --workers 1
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
