# deploy/ubuntu/graffiti-student-worker.service
[Unit]
Description=Graffiti Student Local Queue Worker
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=graffiti
Group=graffiti
WorkingDirectory=/srv/graffiti-student/app
EnvironmentFile=/etc/graffiti-student.env
Environment=PYTHONUNBUFFERED=1
ExecStart=/srv/graffiti-student/venv/bin/python -m deploy.local_worker
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
